# Library/LaunchDaemons/com.apple.corecontextd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.corecontextd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.mobile.keybagd.lock_status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.lock_status</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.corecontextd.xpc</key>
		<true/>
		<key>com.apple.corecontextd.bulletinboard.xpc</key>
		<true/>
		<key>com.apple.private.alloy.tdg-xp-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/libexec/corecontextd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>KeepAlive</key>
	<true/>
	<key>Disabled</key>
	<true/>
</dict>
</plist>
