# Library/LaunchDaemons/com.apple.bbswbypass.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>EnablePressuredExit</key>
	<false/>
	<key>Label</key>
	<string>com.apple.bbswbypass</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.bbswbypass.msgport</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/bbswbypass</string>
	</array>
	<key>BasebandUSBDescriptorEUR</key>
	<dict>
		<key>VendorID</key>
		<real>1478</real>
		<key>ProductID</key>
		<integer>36881</integer>
		<key>Class</key>
		<integer>0</integer>
		<key>Subclass</key>
		<integer>0</integer>
		<key>Protocol</key>
		<integer>0</integer>
		<key>Configuration</key>
		<integer>1</integer>
		<key>CurrentDraw</key>
		<integer>250</integer>
		<key>PollingIntervalMS</key>
		<integer>32</integer>
		<key>Attributes</key>
		<integer>0</integer>
		<key>Label</key>
		<string>IOUSBCompositeDriver</string>
		<key>DIAGInterfaceName</key>
		<string>AppleBB_GNRC1</string>
		<key>DUNInterfaceName</key>
		<string>AppleBB_GNRC2</string>
	</dict>
	<key>BasebandUSBDescriptorICE</key>
	<dict>
		<key>VendorID</key>
		<real>32903</real>
		<key>ProductID</key>
		<integer>2286</integer>
		<key>Class</key>
		<integer>239</integer>
		<key>Subclass</key>
		<integer>2</integer>
		<key>Protocol</key>
		<integer>1</integer>
		<key>Configuration</key>
		<integer>1</integer>
		<key>CurrentDraw</key>
		<integer>100</integer>
		<key>PollingIntervalMS</key>
		<integer>1</integer>
		<key>Attributes</key>
		<integer>224</integer>
		<key>Label</key>
		<string></string>
		<key>DIAGInterfaceName</key>
		<string>CDC ACM (Trace)</string>
		<key>DUNInterfaceName</key>
		<string>CDC ACM (Modem)</string>
	</dict>
	<key>UserName</key>
	<string>root</string>
</dict>
</plist>
