# AppleInternal/Library/LaunchDaemons/com.apple.internal.hello_xpc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.perf.hello_xpc</string>
	<key>Disabled</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.perf.hello_xpc</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/libexec/hello_xpcd</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/tmp/hello_xpcd.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/hello_xpcd.log</string>
</dict>
</plist>
