# AppleInternal/Library/LaunchDaemons/com.apple.mediaservicesbroker.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>UserName</key>
	<string>mobile</string>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>Label</key>
	<string>com.apple.mediaservicesbroker.app</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.agora.server</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/AppleInternal/Applications/Agora.app/Agora</string>
	</array>
	<key>EnableTransactions</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
</dict>
</plist>
