# AppleInternal/Library/LaunchDaemons/com.apple.usbiotestd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.usbiotestd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.appleusbtestdevicecontrol.iokitlaunch</key>
			<dict>
				<key>IOProviderClass</key>
				<string>IOUSBDeviceInterface</string>
				<key>IOPropertyMatch</key>
				<dict>
					<key>USBDeviceFunction</key>
					<string>AppleUSBTestControlInterface</string>
				</dict>
				<key>IOMatchLaunchStream</key>
				<true/>
			</dict>
			<key>com.apple.appleusbtestdevicebulk.iokitlaunch</key>
			<dict>
				<key>IOProviderClass</key>
				<string>IOUSBDeviceInterface</string>
				<key>IOPropertyMatch</key>
				<dict>
					<key>USBDeviceFunction</key>
					<string>AppleUSBTestBulkInterface</string>
				</dict>
				<key>IOMatchLaunchStream</key>
				<true/>
			</dict>
			<key>com.apple.appleusbtestdeviceinterrupt.iokitlaunch</key>
			<dict>
				<key>IOProviderClass</key>
				<string>IOUSBDeviceInterface</string>
				<key>IOPropertyMatch</key>
				<dict>
					<key>USBDeviceFunction</key>
					<string>AppleUSBTestInterruptInterface</string>
				</dict>
				<key>IOMatchLaunchStream</key>
				<true/>
			</dict>
			<key>com.apple.appleusbtestdeviceisoc.iokitlaunch</key>
			<dict>
				<key>IOProviderClass</key>
				<string>IOUSBDeviceInterface</string>
				<key>IOPropertyMatch</key>
				<dict>
					<key>USBDeviceFunction</key>
					<string>AppleUSBTestIsocInterface</string>
				</dict>
				<key>IOMatchLaunchStream</key>
				<true/>
			</dict>
			<key>com.apple.appleusbtestcreateconfiguration.iokitlaunch</key>
			<dict>
				<key>IOProviderClass</key>
				<string>IOUSBDeviceController</string>
				<key>IOPropertyMatch</key>
				<dict>
					<key>EnableTestConfig</key>
					<true/>
				</dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOMatchAll</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>Program</key>
	<string>/usr/local/libexec/usbiotestd</string>
	<key>KeepAlive</key>
	<false/>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
